# config_files/setup.service
[Unit]
Description="My setup service"

[Service]
ExecStart=/home/drb/mysetup.sh

[Install]
WantedBy=multi-user.target
